# docker/celery/celery_supervisord.conf
[program:accounting_worker]
command=/app/docker/celery/worker_up_script
autostart=true
autorestart=true
stderr_logfile=/var/log/celery_err.log
stdout_logfile=/var/log/celery_out.log
